# systemd/timetracker.service
[Unit]
Description=track time

[Service]
Type=oneshot
ExecStart=/usr/local/bin/timetracker-cli
User=timetracker
Group=systemd-journal