# systemd/timetracker.timer
[Unit]
Description=Run timetracker every 5 minutes

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target